# slyd.conf
description "slyd server"
start on mountall
stop on shutdown

chdir /home/portia/portia/slyd
script
    PYTHONPATH=.:./../slybot bin/slyd -p 9001 --root ./dist
end script
respawn
